# samples/laravel-crud-boilerplate/.docker/etc/supervisor/websocket-server.conf
[program:websocket-server]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www/app
command=/usr/bin/laravel-echo-server start
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/websocket-server.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=50
loglevel=warn